# update-notifier.service
[Unit]
Description=Refresh PackageKit cache

[Service]
Type=oneshot
ExecStart=/usr/bin/pkcon refresh -p
