# systemd/update-notifier.service
[Unit]
Description=Refresh PackageKit cache

[Service]
Type=oneshot
ExecStart=/usr/lib/update-notifier/refresh-packagekit-cache
