# arch/update-notifier.install
post_install() {
    echo "Activate update-notifier with"
    echo "    systemctl enable update-notifier.timer"
}